# files/etc/uci-defaults/99-easepi-r1-base
#!/bin/sh
set -e

# Keep image customization minimal by default.
# Add your own UCI customizations under files/etc/uci-defaults/ when needed.

exit 0
